# deploy/systemd/xau-state-discord.service
[Unit]
Description=XAU State Discord Broadcaster
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/quant/calc/calc-xau-engine
Environment=NODE_ENV=production
EnvironmentFile=/quant/calc/config/xau-state-discord.env
ExecStart=/usr/local/bin/node dist/main.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/quant/calc/data/xau-state-discord
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
