# deploy/systemd/xauusd-quant-bot.service
[Unit]
Description=XAUUSD Quantitative Analysis Discord Broadcaster
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/quant/calc/xauusd-quant-bot
Environment=NODE_ENV=production
EnvironmentFile=/quant/calc/config/xauusd-quant-bot.env
ExecStart=/usr/local/bin/node dist/main.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/quant/calc/data/xauusd-quant-bot
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
